# rc/spambox.dat
#!/bin/sh
# chkconfig: 235 99 10
# description: Start or stop the SPAMBOX SPAM Proxy  server
#
### BEGIN INIT INFO
# Provides: SPAMBOX
# Required-Start: $network $syslog
# Required-Stop: $network
# Default-Start: 2 3 5
# Default-Stop: 0 1 6
# Description: Start or stop the SPAMBOX SPAM Proxy  server
### END INIT INFO

# Rev 0.2.1 Lars Troen
# added reload option

# Rev 0.2.0 R.Toth
# Mods to specify the correct startup directory..  
# This needs to be hard-coded for now, until a better way comes up...
# Specifically, the location(s) of the scripts, pid file, config file...
# Mods made for 0.2.0 and new config file format to locate PID file too..

base=/usr/local/spambox
confFile=$base/spambox.cfg
start="$base/start $base"
stop="$base/stop $base"
lockfile=/var/lock/subsys/spambox.lock
pidFile=$base/pid
name='SPAMBOX'

case "$1" in
'start')
	echo "Starting SPAMBOX via: $start"
	$start >/dev/null 2>&1 </dev/null
	# $start 
	RETVAL=$?
	if [ "$RETVAL" = "0" ]; then
		touch $lockfile >/dev/null 2>&1
	fi
	;;
'stop')
	$stop
	RETVAL=$?
	if [ "$RETVAL" = "0" ]; then
		rm -f $lockfile
	fi
	;;
'status')
	pidfile=`grep "^pidfile=" $confFile | sed -e 's/pidfile=//g'`
	if [ "$pidfile" = "" ]; then
		pidfile=$pidFile
	fi
	if [ -s $pidfile ]; then
		pid=`cat $pidfile`
		kill -0 $pid >/dev/null 2>&1
		if [ "$?" = "0" ]; then
			echo "$name (pid $pid) is running"
			RETVAL=0
		else
			echo "$name is stopped"
			RETVAL=1
		fi
	else
		echo "$name is stopped"
		RETVAL=1
	fi
	;;
'restart')
	$stop && $start
	RETVAL=$?
	;;
'reload')
	pidfile=`grep "^pidfile=" $confFile | sed -e 's/pidfile=//g'`
	if [ "$pidfile" = "" ]; then
		pidfile=$pidFile
	fi
	if [ -s $pidfile ]; then
	    pid=`cat $pidfile`
		kill -1 $pid >/dev/null 2>&1
		if [ "$?" = "0" ]; then
			echo "$name (pid $pid) reloading config"
			RETVAL=0
		else
			echo "Unable to reload config. $name is stopped"
			RETVAL=1
		fi
	else
		echo "Unable to reload config. $name is stopped"
		RETVAL=1
	fi
	;;
			
*)
	echo "Usage: $0 { start | stop | status | restart }"
	RETVAL=1
	;;
esac
exit $RETVAL

